# root/etc/uci-defaults/20_99-last-event
#!/bin/sh
# shellcheck disable=SC1091
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :

. /lib/overthebox

otb_save_event "first-setup"
