# root/etc/uci-defaults/20_random
#!/bin/sh

uci -q get system.@rngd[0] > /dev/null && exit 0

uci -q batch <<EOF
add system rngd
set system.@rngd[0].enabled=1
set system.@rngd[0].device=/dev/urandom
commit system
EOF

exit 0
